# bin/systemd/securetea.service
[Unit]
Description=SecureTea

[Service]
ExecStart=/usr/bin/python /usr/local/bin/SecureTea.py
Type=forking
#If not using default pid, make sure you change PIDFile setting
#PIDFile=/var/run/securetea.pid

[Install]
WantedBy=multi-user.target
